# bin/launchd/com.mirador.nightly-smoke.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  ~ macOS launchd unit for bin/nightly-smoke.sh.
  ~
  ~ Install:
  ~   cp bin/launchd/com.mirador.nightly-smoke.plist ~/Library/LaunchAgents/
  ~   launchctl load ~/Library/LaunchAgents/com.mirador.nightly-smoke.plist
  ~
  ~ Verify it loaded:
  ~   launchctl list | grep mirador
  ~
  ~ Uninstall:
  ~   launchctl unload ~/Library/LaunchAgents/com.mirador.nightly-smoke.plist
  ~   rm ~/Library/LaunchAgents/com.mirador.nightly-smoke.plist
  ~
  ~ Edit the WorkingDirectory + Program paths if your repo lives elsewhere.
  ~ stdout / stderr land in /tmp/com.mirador.nightly-smoke.{out,err}.log.
  -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mirador.nightly-smoke</string>

    <key>WorkingDirectory</key>
    <string>/Users/benoitbesson/dev/workspace-modern/mirador-service</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-l</string>
        <string>-c</string>
        <string>./bin/nightly-smoke.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/com.mirador.nightly-smoke.out.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/com.mirador.nightly-smoke.err.log</string>
</dict>
</plist>
